# eink-cal.timer
[Unit]
Description = run eink-cal.py hourly

[Timer]
OnCalendar = hourly

[Install]
WantedBy = timers.target
